# config/launchd/com.moltbot.security-audit.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.moltbot.security-audit.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/ss/moltbot/scripts/run-security-audit.sh</string>
        <string>--profile</string>
        <string>default</string>
        <string>--alert</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/moltbot-security-audit-daily.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/moltbot-security-audit-daily.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
